# com.nhnews.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nhnews.daemon</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/splendor/antigravity/nh_news/venv/bin/python3</string>
        <string>/Users/splendor/antigravity/nh_news/main.py</string>
        <string>--loop</string>
        <string>--interval</string>
        <string>10</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/splendor/antigravity/nh_news</string>
    <key>StandardOutPath</key>
    <string>/Users/splendor/antigravity/nh_news/nh_news.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/splendor/antigravity/nh_news/nh_news.log</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
</dict>
</plist>
